# serge-zing/crontab
*/5 * * * * serge sync $SERGE_SETTINGS >/var/log/serge.log 2>&1
