# serge/crontab
*/5 * * * * serge sync /etc/serge/serge.conf >/var/log/serge.log 2>&1
